# debian/open-webui-app.service
[Unit]
Description=Open WebUI Docker Container
# Only start after Docker is ready
After=docker.service
After=open-webui-ollama.service
After=open-webui-postgresql.service
Requires=docker.service
Requires=open-webui-ollama.service
Requires=open-webui-postgresql.service
# Wait for network to be available
After=network-online.target

[Service]
Type=simple

# Environment files
EnvironmentFile=/etc/default/open_webui
EnvironmentFile=/etc/open-webui/open-webui-app.conf
EnvironmentFile=-/run/open-webui/open-webui-app.conf

# Remove any existing container with this name before starting
ExecStartPre=-/usr/bin/docker rm -f owebui_app
# Check if Docker image is available
ExecStartPre=/usr/lib/open-webui-tools/owebui_app_env_cond \
    /etc/open-webui/open-webui-app.conf \
    /etc/open-webui/open-webui-local.conf \
    /etc/default/open_webui \
    /run/open-webui/open-webui-app.conf
# Generate dynamic environment file
ExecStartPre=-/usr/bin/mkdir -p /run/open-webui
ExecStartPre=/usr/lib/open-webui-tools/owebui_app_env_pre \
    /etc/open-webui/open-webui-ollama.conf \
    /etc/open-webui/open-webui-postgresql.conf \
    /etc/open-webui/open-webui-app.conf \
    /etc/default/open_webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-app.conf

# Start the container in the foreground (no -d flag)
ExecStart=/usr/bin/docker run \
  --name owebui_app \
  --rm \
  $PORTMAP_CMD $PORTMAP_ARG \
  -v ${OWEBUI_APP}:/app/backend/data \
  $OWEBUI_DEVICE \
  $OWEBUI_ENV \
  ${OWEBUI_IMAGE}

# Gracefully stop the container
ExecStop=/usr/bin/docker stop -t 30 owebui_app

# Timeout for start operation
TimeoutStartSec=60
# Timeout for stop operation
TimeoutStopSec=35

# Always restart on failure
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
