# debian/open-webui-docling.service
[Unit]
Description=Ollama Docker Container
# Only start after Docker is ready
After=docker.service
Requires=docker.service
# Wait for network to be available
After=network-online.target

[Service]
Type=simple

# Environment files
EnvironmentFile=/etc/open-webui/open-webui-docling.conf
EnvironmentFile=/etc/default/open_webui
EnvironmentFile=-/etc/open-webui/open-webui-local.conf
EnvironmentFile=-/run/open-webui/open-webui-docling.conf

# Remove any existing container with this name before starting
ExecStartPre=-/usr/bin/docker rm -f owebui_docling
# Check if Docker image is available
ExecStartPre=/usr/lib/open-webui-tools/owebui_docling_env_cond \
    /etc/open-webui/open-webui-docling.conf \
    /etc/default/open_webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-docling.conf
# Generate dynamic environment file
ExecStartPre=-/usr/bin/mkdir -p /run/open-webui
ExecStartPre=/usr/lib/open-webui-tools/owebui_docling_env_pre \
    /etc/open-webui/open-webui-docling.conf \
    /etc/default/open_webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-docling.conf

ExecStopPost=/usr/lib/open-webui-tools/owebui_docling_env_stop \
    /etc/open-webui/open-webui-docling.conf \
    /etc/default/open_webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-docling.conf

# Start the container in the foreground (no -d flag)
ExecStart=/usr/bin/docker run \
  --name owebui_docling \
  --rm \
  $PORTMAP_CMD $PORTMAP_ARG \
  $DOCLING_DEVICE \
  $DOCLING_ENV \
  ${DOCLING_IMAGE}

# Gracefully stop the container
ExecStop=/usr/bin/docker stop -t 30 owebui_docling

# Timeout for start operation
TimeoutStartSec=60
# Timeout for stop operation
TimeoutStopSec=35

# Always restart on failure
Restart=always
RestartSec=45

[Install]
WantedBy=multi-user.target
