# debian/open-webui-ollama.service
[Unit]
Description=Ollama Docker Container
# Only start after Docker is ready
After=docker.service
Requires=docker.service
# Wait for network to be available
After=network-online.target

[Service]
Type=simple

# Environment files
EnvironmentFile=/etc/open-webui/open-webui-ollama.conf
EnvironmentFile=/etc/default/open-webui
EnvironmentFile=-/etc/open-webui/open-webui-local.conf
EnvironmentFile=-/run/open-webui/open-webui-ollama.conf

# Remove any existing container with this name before starting
ExecStartPre=-/usr/bin/docker rm -f owebui_ollama
# Check if Docker image is available
ExecStartPre=/usr/lib/open-webui-tools/owebui_ollama_env_cond \
    /etc/open-webui/open-webui-ollama.conf \
    /etc/default/open-webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-ollama.conf
# Generate dynamic environment file
ExecStartPre=-/usr/bin/mkdir -p /run/open-webui
ExecStartPre=/usr/lib/open-webui-tools/owebui_ollama_env_pre \
    /etc/open-webui/open-webui-ollama.conf \
    /etc/default/open-webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-ollama.conf

# Post process
ExecStopPost=/usr/lib/open-webui-tools/owebui_ollama_env_stop \
    /etc/open-webui/open-webui-ollama.conf \
    /etc/default/open-webui \
    /etc/open-webui/open-webui-local.conf \
    /run/open-webui/open-webui-ollama.conf
ExecStopPost=-/usr/bin/docker rm -f owebui_ollama

# Start the container in the foreground (no -d flag)
ExecStart=/usr/bin/docker run \
  --name owebui_ollama \
  --rm \
  $OWEBUI_DOCKER_OPTS \
  $PORTMAP_CMD $PORTMAP_ARG \
  -v ${OWEBUI_OLLAMA}:/root/.ollama \
  $OLLAMA_DEVICE \
  $OLLAMA_ENV \
  ${OLLAMA_IMAGE}

# Gracefully stop the container
ExecStop=/usr/bin/docker stop -t 30 owebui_ollama

# Timeout for start operation
TimeoutStartSec=60
# Timeout for stop operation
TimeoutStopSec=35

# Always restart on failure
Restart=always
RestartSec=45

[Install]
WantedBy=multi-user.target
